# 99-monitor-hotplug.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/bin/bash /PATH/TO/monitor-hotplug.sh"
